# addon/rootfs/etc/s6-overlay/s6-rc.d/postgres-setup/up
#!/command/execlineb -P
/usr/local/bin/scryvex-postgres-setup
