# install/docker/fs/etc/s6-overlay/s6-rc.d/setup/up
/etc/s6-overlay/scripts/setup.sh && /etc/s6-overlay/scripts/seed-plugins.sh